# .config/systemd/user/swaync.service.d/override.conf
[Unit]
ConditionEnvironment=

[Service]
PassEnvironment=WAYLAND_DISPLAY XDG_CURRENT_DESKTOP
