# .config/systemd/user/user/swaybg.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/swaybg -m fill -i "/home/jraj/Pictures/wallies/IMG_0130.jpg"

