# .config/systemd/user/swaybg.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/swaybg -m fill -i "/home/jraj/Pictures/wally/wallhaven-5g9rl7_3840x2160.png"

